# roles/workstation/templates/default/etc/profile.d/ghq.sh.j2
#!/bin/bash -eu
export PATH=$PATH:/opt/ghq
