# roles/workstation-system/templates/default/etc/profile.d/go.sh.j2
#!/bin/bash -eu
export PATH=$PATH:/opt/go/bin
